# scripts/systemd-user/openclaw-health.service
[Unit]
Description=OpenClaw health check (script-only, notify on critical)

[Service]
Type=oneshot
Environment=HOME=%h
ExecStart=/home/teaingtit/projects/openclaw/scripts/ops/health-check.sh --format json --notify-on-critical
SuccessExitStatus=1 2
