# scripts/systemd-user/openclaw-sot-sync.service
[Unit]
Description=Regenerate SYSTEM_INDEX.md agent table from openclaw.json

[Service]
Type=oneshot
Environment=HOME=%h
WorkingDirectory=/home/teaingtit/projects/openclaw
ExecStart=/home/teaingtit/projects/openclaw/scripts/ops/gen-agent-index.sh
SuccessExitStatus=0 1
